# docker/container-data/etc/cron.d/project-crontab
## MAILING

# [Each 1 min] Send E-Mails
SEND_EMAILS_LOG_FILE_PATH=/var/log/cron-send-emails.log
* * * * * root cd /var/www/html/ &&  date >> $SEND_EMAILS_LOG_FILE_PATH && php bin/console voltigo:email:send-emails >> $SEND_EMAILS_LOG_FILE_PATH 2>&1

# [Each 1 min] Check E-Mails sending status
CHECK_EMAILS_SENDING_STATUS_LOG_FILE_PATH=/var/log/cron-check-emails-sending-status.log
* * * * * root cd /var/www/html/ &&  date >> $CHECK_EMAILS_SENDING_STATUS_LOG_FILE_PATH && php bin/console voltigo:email:check-status >> $CHECK_EMAILS_SENDING_STATUS_LOG_FILE_PATH 2>&1
